# divertscan-capture.service
[Unit]
Description=DivertScan Scale Capture Service
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/divertscan
ExecStart=/usr/bin/python3 /home/pi/divertscan/scale_capture.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Watchdog - restart if process hangs
WatchdogSec=120

# Environment
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
