# scale_capture.service
[Unit]
Description=DivertScan Scale Capture
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi
ExecStart=/usr/bin/python3 /home/pi/scale_capture.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
